# farm_irrigation/farm_irrigation.service
# systemd服务配置文件 - 农场级灌溉设计算法AI
# 放置位置: /etc/systemd/system/farm_irrigation.service
# 启用服务: sudo systemctl enable farm_irrigation
# 启动服务: sudo systemctl start farm_irrigation

[Unit]
Description=Farm Irrigation Design Algorithm AI Web Application
After=network.target
Wants=network.target

[Service]
# 服务类型
Type=exec

# 用户和组（建议创建专用用户）
User=farm_app
Group=farm_app

# 工作目录
WorkingDirectory=/opt/farm_irrigation

# 环境变量
Environment=PATH=/opt/farm_irrigation/venv/bin
Environment=FLASK_APP=farm_irrigation/web_farm_irrigation_modified.py
Environment=FLASK_ENV=production
Environment=PYTHONPATH=/opt/farm_irrigation

# 启动命令
ExecStart=/opt/farm_irrigation/venv/bin/gunicorn --config /opt/farm_irrigation/gunicorn.conf.py web_farm_irrigation_modified:app

# 重启策略
Restart=always
RestartSec=10

# 进程管理
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=farm_irrigation

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/farm_irrigation
ReadWritePaths=/var/log/farm_irrigation
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target